# services/shout/run
#!/bin/bash
exec /sbin/setuser abc shout --home /config
